# debian/cron.d
#
# cron jobs for xivo-purge-db
#

25 3 * * * root /usr/bin/xivo-purge-db
